# web_scanner.service.example
# systemd service example for web_scanner
# Copy to: /etc/systemd/system/web_scanner.service
# Then: sudo systemctl daemon-reload && sudo systemctl enable --now web_scanner

[Unit]
Description=OpenRouter Scanner Web Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/hajimi-king-main
ExecStart=/usr/bin/python3 web_scanner.py
# Optional environment overrides (uncomment as needed)
# Environment=HOST=0.0.0.0
# Environment=PORT=5000
# Environment=DEBUG=false
# Environment=SCANNER_DB_PATH=/var/lib/web_scanner/scanner.db
# Environment=ALLOWED_ORIGINS=https://your-admin.example.com
# Environment=OPENROUTER_VALIDATION_ENDPOINT=https://openrouter.ai/api/v1/auth/key
# Environment=OPENROUTER_TIMEOUT=10
# Environment=ADMIN_BEARER=your-secret-token
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

